# config/war3bot.service
[Unit]
Description=War3Bot Warcraft III Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/War3Bot/build
ExecStart=/root/War3Bot/build/war3bot -p 6112
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target